# systemd/kortexa-ai-tts.service
[Unit]
Description=Kortexa TTS Server
After=network.target

StartLimitIntervalSec=120
StartLimitBurst=3

[Service]
Type=simple

User=francip
Group=francip

WorkingDirectory=/home/francip/src/tts.server

Environment=PATH=/home/francip/.local/bin:/home/francip/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/home/francip/src/tts.server

ExecStart=/home/francip/src/tts.server/run.sh

Restart=always
RestartSec=15

TimeoutStopSec=60

ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/home/francip/src/tts.server /home/francip/.cache /home/francip/.config /home/francip/data
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
