# launchd/ai.kortexa.tts.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/Users/francip/.local/bin:/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin</string>
		<key>PYTHONPATH</key>
		<string>/Users/francip/src/tts.server</string>
		<key>PYTHONUNBUFFERED</key>
		<string>1</string>
	</dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>ai.kortexa.tts.server</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/francip/src/tts.server/launchd/kortexa-tts-server.sh</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/francip/Library/Logs/tts.log</string>
	<key>StandardOutPath</key>
	<string>/Users/francip/Library/Logs/tts.log</string>
	<key>ThrottleInterval</key>
	<integer>15</integer>
	<key>WorkingDirectory</key>
	<string>/Users/francip/src/tts.server</string>
</dict>
</plist>
