# launchd/ai.kortexa.tts.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.kortexa.tts</string>

    <key>WorkingDirectory</key>
    <string>/Users/francip/src/tts.server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/francip/src/tts.server/launchd/kortexa-tts.sh</string>
        <string>--dev</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/francip/.local/bin:/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
        <key>PYTHONPATH</key>
        <string>/Users/francip/src/tts.server</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>15</integer>

    <key>StandardOutPath</key>
    <string>/Users/francip/Library/Logs/tts.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/francip/Library/Logs/tts.log</string>
</dict>
</plist>
